# apps/mercury-agent/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx.out.log
stderr_logfile=/var/log/supervisor/nginx.err.log

[program:ttyd]
command=/usr/bin/ttyd -p 7681 -i 127.0.0.1 -W -t titleFixed=Mercury -t fontSize=13 -T xterm-256color /app/run-mercury.sh
directory=/data
autorestart=true
priority=20
environment=TERM="xterm-256color",MERCURY_HOME="/data",MERCURY_BIN="/usr/local/bin/mercury"
stdout_logfile=/var/log/supervisor/ttyd.out.log
stderr_logfile=/var/log/supervisor/ttyd.err.log

[program:mercury-daemon]
command=/app/run-mercury-daemon.sh
directory=/data
autorestart=true
startsecs=10
startretries=10000
priority=30
environment=MERCURY_HOME="/data",HOME="/root",PATH="/usr/local/bin:/usr/bin:/bin"
stdout_logfile=/var/log/supervisor/mercury-daemon.out.log
stderr_logfile=/var/log/supervisor/mercury-daemon.err.log
